# nas/caddy-prepare.service
[Unit]
Description=Prepare persistent Caddy state
RequiresMountsFor=/var/srv/ssd
Before=caddy.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/infrastructure/caddy-prepare
RemainAfterExit=yes
