# nas/nfs-firewall-prepare.service
[Unit]
Description=Restrict NFSv4 to the reserved VM addresses
Wants=network-online.target
After=network-online.target nftables.service
Requires=nftables.service
Before=nfs-server.service
RequiresMountsFor=/var/srv/docs

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/infrastructure/nfs-firewall-prepare

[Install]
WantedBy=multi-user.target
