# ips/clamav-freshclam.service.d/10-prepare.conf
[Unit]
Wants=network-online.target clamav-prepare.service
After=network-online.target clamav-prepare.service
