# workstation/fapolicyd-flatpak-sync.path
[Unit]
Description=Watch system Flatpak changes for fapolicyd trust sync
ConditionPathIsDirectory=/var/lib/flatpak

[Path]
PathChanged=/var/lib/flatpak/.changed
PathChanged=/var/lib/flatpak/app
PathChanged=/var/lib/flatpak/runtime
Unit=fapolicyd-flatpak-sync.service

[Install]
WantedBy=multi-user.target
