# ips/suricata.service.d/10-inline.conf
[Unit]
Wants=network-online.target suricata-prepare.service
After=network-online.target suricata-prepare.service

[Service]
ExecStartPre=/usr/bin/test -s /var/lib/suricata/rules/suricata.rules
