# ips/suricata.service.d/10-inline.conf
[Unit]
Wants=network-online.target suricata-prepare.service
After=network-online.target suricata-prepare.service
Requires=suricata-prepare.service
