# ips/suricata.service.d/10-inline.conf
[Unit]
Wants=network-online.target
After=network-online.target

[Service]
ExecStartPre=/usr/bin/test -s /var/lib/suricata/rules/suricata.rules
