# workstation/gssproxy.service.d/10-bootcontainer.conf
[Unit]
Wants=gssproxy-prepare.service
After=gssproxy-prepare.service

[Service]
StateDirectory=
StateDirectoryMode=
