# ips/suricata-prepare.service
[Unit]
Description=Prepare Suricata rules before enabling inline filtering
Wants=network-online.target
After=network-online.target
Before=suricata.service
ConditionPathExists=!/var/lib/suricata/rules/suricata.rules

[Service]
Type=oneshot
ExecStart=/usr/libexec/bootcontainer/suricata-prepare-rules

[Install]
WantedBy=multi-user.target
