# images/workstation/home-ownership.service
[Unit]
Description=Ensure local user home directory ownership and mode
After=local-fs.target systemd-sysusers.service systemd-tmpfiles-setup.service
Before=systemd-user-sessions.service apply-skel.service display-manager.service sshd.service
Wants=systemd-tmpfiles-setup.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/infrastructure/home-ownership

[Install]
WantedBy=multi-user.target
