# ips/clamd@scan.service.d/10-prepare.conf
[Unit]
Wants=clamav-prepare.service
After=clamav-prepare.service
