# nas/containers/nextcloud.container
[Unit]
Description=Nextcloud
ConditionPathExists=/var/lib/containers/nextcloud

[Container]
Image=docker.io/library/nextcloud:apache
Environment=SQLITE_DATABASE=nextcloud
Environment=NEXTCLOUD_TRUSTED_DOMAINS=nick.nas
Environment=OVERWRITEPROTOCOL=https
Environment=OVERWRITECLIURL=https://nick.nas/nextcloud
Environment=OVERWRITEWEBROOT=/nextcloud
Memory=1G
PidsLimit=512
PublishPort=127.0.0.1:8081:80/tcp
Volume=/var/lib/containers/nextcloud:/var/www/html:Z
Volume=/var/srv/docs/shared:/var/srv/docs/shared:rw,Z,nodev,nosuid,noexec
AutoUpdate=registry

[Service]
Restart=always
RestartSec=5
TimeoutStartSec=900
TimeoutStopSec=180
MemoryHigh=1536M
MemoryMax=2G
CPUQuota=200%
IOWeight=50

[Install]
WantedBy=multi-user.target
