# nas/containers/paperless-tika.container
[Unit]
Description=Paperless-ngx Apache Tika document parser

[Container]
Image=docker.io/apache/tika:latest
ContainerName=paperless-tika
Network=paperless.network
NetworkAlias=tika
Memory=1G
PidsLimit=512
AutoUpdate=registry

[Service]
Restart=always
RestartSec=5
TimeoutStartSec=900
TimeoutStopSec=120
MemoryHigh=1G
MemoryMax=1536M
CPUQuota=200%
IOWeight=50

[Install]
WantedBy=multi-user.target
