# nas/containers/paperless-db.container
[Unit]
Description=Paperless-ngx PostgreSQL database
Requires=paperless-prepare.service
After=paperless-prepare.service

[Container]
Image=docker.io/library/postgres:18
ContainerName=paperless-db
Network=paperless.network
NetworkAlias=db
Environment=POSTGRES_DB=paperless
Environment=POSTGRES_USER=paperless
Environment=POSTGRES_PASSWORD_FILE=/run/secrets/db-password
Memory=1G
PidsLimit=512
Volume=/var/lib/containers/paperless/postgres:/var/lib/postgresql:Z
Volume=/var/lib/containers/secrets/paperless/db-password:/run/secrets/db-password:ro,z
AutoUpdate=registry

[Service]
Restart=always
RestartSec=5
TimeoutStartSec=900
TimeoutStopSec=180
MemoryHigh=1G
MemoryMax=1536M
CPUQuota=150%
IOWeight=50

[Install]
WantedBy=multi-user.target
